# packages/sys-power/apcupsd/files/systemd/apcupsd.service
[Unit]
Description=APC UPS monitoring daemon
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/sbin/apcupsd -b

[Install]
WantedBy=multi-user.target
